# website/durable-request-web.service
[Unit]
Description=Durable Request Product Website (Internal LAN)
After=network.target

[Service]
Type=simple
User=albert
WorkingDirectory=/home/albert/durable-request/website
ExecStart=/home/albert/durable-request/website/serve.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
